# brainrunner.service
[Unit]
Description=brainrunner — autonomous GitHub issue agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@USER@
Environment="PATH=@HOME@/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/usr/local/bin/brainrunner --config @HOME@/.config/brainrunner/config.toml
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=brainrunner

[Install]
WantedBy=multi-user.target
